# IoT/hssi_iot.timer
[Unit]
Description=HSSI IoT ever minite

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target

