# IoT/hssi_iot.service
[Unit]
Description=HSSI IoT daemon

[Service]
Type=oneshot
ExecStart=/home/root/IoT/iot.sh

[Install]
WantedBy=multi-user.target

